# conf/supervisor/supervisor.d/initdb.ini
[program:initdb]
command=/usr/bin/initdb.sh
priority=9
