# conf/supervisor/supervisor.d/leanote.ini
[program:leanote]
command=/usr/leanote/bin/run.sh
priority=7
