# conf/supervisor/supervisor.d/mongo.ini
[program:mongod]
command=mongod --dbpath /data/mongodb --logpath /data/log/mongodb/mongo.log --fork
priority=5
